# jeedom/jeedom.conf
[program:jeedom]
command=node jeedomold.js
directory=/home/git/seaconnect/jeedom
environment=NODE_ENV=production,NODE_PATH=/home/git/seaconnect/jeedom
stdout_logfile=/var/log/jeedom.log
